# main/dehydrated/dehydrated.timer
[Unit]
Description=Check and renew ACME TLS certificates

[Timer]
OnCalendar=daily
RandomizedDelaySec=7200

[Install]
WantedBy=timers.target
